# .config/autostart/login.sh.desktop
[Desktop Entry]
Exec=/home/matej/.scripts/login.sh
Icon=dialog-scripts
Name=login.sh
Type=Application
X-KDE-AutostartScript=true
